# php-websocket/supervisord.d/laravel-websocket.conf
[program:laravel-websocket]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/html/artisan websockets:serve
autostart=true
autorestart=true
numprocs=8
user=laravel
stdout_logfile=/var/www/html/storage/logs/websockets.log
