# php-websocket/supervisord.d/laravel-websocket.conf
[program:laravel-websocket]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/html/artisan websockets:serve --port 6001
autostart=true
autorestart=true
numprocs=8
user=laravel
